# systemd-services/.config/systemd/user/kmonad.service
[Unit]
Description=KMonad

[Service]
ExecStart=kmonad %h/.config/kmonad/main.kbd

[Install]
WantedBy=default.target