# systemd-services/.config/systemd/user/ydotool.service
[Unit]
Description=ydotool

[Service]
ExecStart=ydotoold

[Install]
WantedBy=default.target