# systemd/.config/systemd/user/kmonad.service
[Unit]
Description=KMonad

[Service]
ExecStart=%h/.guix-profile/bin/kmonad %h/.config/kmonad/main.kbd

[Install]
WantedBy=default.target